# device/deploy/s3-uploader.service
[Unit]
Description=Smart Camera S3 Uploader
After=smart-camera.service
Requires=smart-camera.service

[Service]
Type=simple
ExecStart=/opt/smart-camera/ai/run_uploader.sh
WorkingDirectory=/opt/smart-camera
Restart=always
RestartSec=10
Environment=AI_SUMMARY_CONFIG_PATH=/etc/smart-camera/config.ini
Environment=AWS_DEFAULT_REGION=ap-southeast-1
Environment=DEVICE_ID=raspi-camera-01

[Install]
WantedBy=multi-user.target
